# PAN_Pi4/rc.local
#!/bin/sh -e
#
# rc.local
sudo bash /interface_restore.sh
exit 0